# ThisFolderBkp/BKP_Aug_26_2025/apoorva_driver_autosetup/etc/systemd/system/apoorva1553.service
[Unit]
Description=Apoorva PCI MIL-STD-1553B Driver
After=network.target

[Service]
Type=oneshot
ExecStart=/sbin/modprobe ApoorvaPciDriver
ExecStartPost=/bin/mknod -m 666 /dev/APOORVA1553 c 236 0
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
